# docker/debian/11/service/usr/lib/systemd/system/wecomsdk.service
[Unit]
Description=Wecom Sdk Api
After=network.target

[Service]
Type=simple
User=odoo
Group=odoo
WorkingDirectory=/code
ExecStart=python3 uvicorn app.main:app  --port 8000
Restart=always


[Install]
WantedBy=multi-user.target